# systemd/user/ydotool.socket
[Unit]
Description=ydotool socket (user)

[Socket]
ListenStream=%t/ydotoold/ydotoold.socket
SocketMode=0600

[Install]
WantedBy=sockets.target
